# files/crontab
* * * * * echo aptible/docker-cron-example >> /var/log/cron.log 2>&1
